# pi-scale-client/scale-client.service
[Unit]
Description=Farm Pi Client (Scale + Barcode Scanner → Server)
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=stepan
SupplementaryGroups=input
WorkingDirectory=/home/stepan/pi-scale-client
ExecStart=/home/stepan/pi-scale-client/venv/bin/python pi_client.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
# Watchdog: если скрипт не шлёт heartbeat 30 сек — считаем зависшим, убиваем
WatchdogSec=30

[Install]
WantedBy=multi-user.target
